# launchd/com.altera.activity-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.altera.activity-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jimmy.leimonitis/.nvm/versions/node/v24.11.0/bin/node</string>
        <string>/Users/jimmy.leimonitis/GitHub/apac-intelligence-v2/scripts/sync-excel-activities.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jimmy.leimonitis/GitHub/apac-intelligence-v2</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/activity-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/activity-sync-error.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jimmy.leimonitis/.nvm/versions/node/v24.11.0/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
